# dockerfiles/nginx-supervisor.ini
[supervisord]
nodaemon=true

[program:nginx]
command = /usr/sbin/nginx
user = root
autostart = true

[program:php-fpm]
command = /usr/sbin/php-fpm7.0 -c /etc/php/7.0/fpm/php-fpm.conf
user = root
autostart = true

[program:cron]
command = /usr/sbin/cron